# jupyter/supervisor-jupyter.conf
[program:juputer]
command=/bin/bash -c /usr/local/bin/jupyter.sh
user=admin
