# node/supervisor-node.conf
[program:node]
command=slurmd -N $(hostname)
