# base/supervisor-base.conf
[program:sed_slurm_conf]
command=/bin/bash -c /etc/slurm-llnl/sed_slurm_conf.sh
autorestart=false
startretries=0

[program:munge]
command=munge
